# profile.d/aliases.sh
# include .aliases if it exists
if [ -f "$HOME/.aliases" ] ; then
    source "$HOME/.aliases"
fi
