# profile.d/mysql-client.sh
export PATH="$HOMEBREW_PREFIX/opt/mysql-client/bin:$PATH"
